# src/pcscd_acsccid.rules
# udev rules for pcscd and ACS CCID readers

# If not adding the device, go away
ACTION!="add", GOTO="pcscd_acsccid_rules_end"

# ACS ACR38U-CCID
SUBSYSTEMS=="usb", ATTRS{idVendor}=="072f", ATTRS{idProduct}=="90cc", RUN+="/usr/sbin/pcscd --hotplug"
# ACS ACR122U
SUBSYSTEMS=="usb", ATTRS{idVendor}=="072f", ATTRS{idProduct}=="2200", RUN+="/usr/sbin/pcscd --hotplug"

# non CCID generic (InterfaceClass: 0x00)
# ACS ACR88U
SUBSYSTEMS=="usb", ATTRS{idVendor}=="072f", ATTRS{idProduct}=="2011", RUN+="/usr/sbin/pcscd --hotplug"
# ACS ACR128U
SUBSYSTEMS=="usb", ATTRS{idVendor}=="072f", ATTRS{idProduct}=="2100", RUN+="/usr/sbin/pcscd --hotplug"

# All done
LABEL="pcscd_acsccid_rules_end"
